# install/linux/check.cron
#!/bin/bash

PID=`pgrep -u recd recd`

if [ -z "${PID}" ]
then
 /home/recd/bin/recd >> /dev/null 
fi

